# systemd/tacplus_ng-reload.path
[Unit]
Description=Watch for /opt/tacacsgui/reload_tacacs_cfg to trigger myapp SIGHUP

[Path]
PathExists=/opt/tacacsgui/reload_tacacs_cfg
Unit=tacplus_ng-reload.service

[Install]
WantedBy=multi-user.target
